# docker/magento/conf/etc/supervisor.d/sync.conf
[program:sync]
command=/bin/bash -c "gosu "www-data" /opt/docker/sync/watch.sh"
startsecs=0
autostart = true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
